# peer_httpserver/supervisord.conf
[supervisord]
nodaemon=true

[program:server]
directory=/
command=node src/index.js
